# lab7/roles/mysql_exporter/templates/mysqld_exporter.service.j2
[Unit]
Description=Node Exporter

[Service]
User=root
ExecStart=/opt/mysqld_exporter/mysqld_exporter

[Install]
WantedBy=multi-user.target
