# EXAM/roles/custom_exporter/templates/custom_exporter.service.j2
[Unit]
Description=Ping Exporter

[Service]
User=root
ExecStart=/usr/local/bin/custom_exporter

[Install]
WantedBy=multi-user.target
